# playbooks/roles/add_restic_server/templates/rest-server.service.j2
[Unit]
Description=rest-server backend for restic backups
Documentation=https://github.com/restic/rest-server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ add_restic_server_user }}
Group={{ add_restic_server_group }}
WorkingDirectory={{ add_restic_server_data_dir }}
UMask=0077
{% set add_restic_server_exec_args = [
  add_restic_server_binary_path,
  '--listen', add_restic_server_listen_address,
  '--path', add_restic_server_data_dir,
  '--log', add_restic_server_log_target
] %}
{% if add_restic_server_append_only | bool %}{% set _ = add_restic_server_exec_args.append('--append-only') %}{% endif %}
{% if add_restic_server_private_repos | bool %}{% set _ = add_restic_server_exec_args.append('--private-repos') %}{% endif %}
{% if add_restic_server_group_accessible_repos | bool %}{% set _ = add_restic_server_exec_args.append('--group-accessible-repos') %}{% endif %}
{% if add_restic_server_prometheus | bool %}{% set _ = add_restic_server_exec_args.append('--prometheus') %}{% endif %}
{% if add_restic_server_prometheus_no_auth | bool %}{% set _ = add_restic_server_exec_args.append('--prometheus-no-auth') %}{% endif %}
{% if add_restic_server_no_verify_upload | bool %}{% set _ = add_restic_server_exec_args.append('--no-verify-upload') %}{% endif %}
{% set _ = add_restic_server_exec_args.append('--htpasswd-file') %}
{% set _ = add_restic_server_exec_args.append(add_restic_server_htpasswd_file) %}
{% if add_restic_server_tls | bool %}{% set _ = add_restic_server_exec_args.extend(['--tls', '--tls-cert', add_restic_server_tls_cert_file, '--tls-key', add_restic_server_tls_key_file, '--tls-min-ver', add_restic_server_tls_min_ver]) %}{% endif %}
{% for flag in add_restic_server_extra_flags %}{% set _ = add_restic_server_exec_args.append(flag) %}{% endfor %}
ExecStart={{ add_restic_server_exec_args | join(' ') }}
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths={{ add_restic_server_data_dir }}

[Install]
WantedBy=multi-user.target
