# playbooks/roles/unbound_setup/templates/unbound-refresh-blocklist.service.j2
[Unit]
Description=Refresh Unbound blocklist configuration
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart={{ unbound_blocklist_refresh_script_path }}
