# playbooks/roles/add_unbound/templates/unbound-refresh-blocklist.timer.j2
[Unit]
Description=Refresh Unbound blocklist configuration on schedule

[Timer]
OnCalendar={{ add_unbound_blocklist_refresh_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target

